# vagrant_files/etc/profile.d/env_vars.sh
#!/bin/sh
export GRADLE_USER_HOME=/vagrant/.gradle
export JAVA_HOME=/usr/lib/jvm/java-1.7.0-openjdk-amd64
